# hotplug.d/iface/99-iptv-wan
if [ "${ACTION}" = "ifup" ] && [ "${INTERFACE}" = "iptv-wan" ]; then

    logger -t hotplug "Interface: '${INTERFACE}' ${ACTION}, adding custom routes..."

fi
