# hotplug.d/iface/99-iptv-wan
if [ "${ACTION}" = "ifup" ] && [ "${INTERFACE}" = "iptv-wan" ]; then
    logger -t hotplug "Interface: '${INTERFACE}' ${ACTION}, adding custom routes..."
    ip route add 61.167.60.0/24 via
fi
